# crons/plists/com.claude.gather.exa.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <!-- Label must match the filename (minus .plist) -->
    <key>Label</key>
    <string>com.claude.gather.exa</string>

    <!-- The script to run -->
    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <!-- __SCRIPT_DIR__ is replaced by manage.sh install -->
        <string>__SCRIPT_DIR__/gather/exa-scan.py</string>
    </array>

    <!-- Schedule: 4:30 AM daily -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>4</integer>
        <key>Minute</key>
        <integer>30</integer>
    </dict>

    <!-- Environment variables -->
    <key>EnvironmentVariables</key>
    <dict>
        <!-- Set your Exa API key here or in your shell profile -->
        <key>EXA_API_KEY</key>
        <string>YOUR_EXA_API_KEY_HERE</string>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin</string>
    </dict>

    <!-- Logging -->
    <key>StandardOutPath</key>
    <string>/tmp/claude-crons-exa-stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/claude-crons-exa-stderr.log</string>

    <!-- Run if machine was asleep at scheduled time -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>4</integer>
        <key>Minute</key>
        <integer>30</integer>
    </dict>

    <!-- Working directory -->
    <key>WorkingDirectory</key>
    <string>__SCRIPT_DIR__</string>

    <!-- Don't keep running — one-shot -->
    <key>KeepAlive</key>
    <false/>
</dict>
</plist>
